# root/etc/services.d/unmanic/run
#!/usr/bin/with-contenv bash

echo "Starting UnManic."
cd /app
echo $PWD
exec \
s6-setuidgid abc python3 /app/service.py